# build/com.isp-snitch.monitor.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.isp-snitch.monitor</string>
    
    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/isp-snitch</string>
        <string>start</string>
        <string>--background</string>
    </array>
    
    <key>RunAtLoad</key>
    <true/>
    
    <key>KeepAlive</key>
    <true/>
    
    <key>StandardOutPath</key>
    <string>/usr/local/var/log/isp-snitch/out.log</string>
    
    <key>StandardErrorPath</key>
    <string>/usr/local/var/log/isp-snitch/error.log</string>
    
    <key>WorkingDirectory</key>
    <string>/usr/local/var/isp-snitch</string>
    
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/usr/local/bin:/usr/bin:/bin</string>
        <key>ISP_SNITCH_CONFIG</key>
        <string>/usr/local/etc/isp-snitch/config.json</string>
        <key>ISP_SNITCH_DATA</key>
        <string>/usr/local/var/isp-snitch/data</string>
    </dict>
    
    <key>ThrottleInterval</key>
    <integer>10</integer>
    
    <key>ProcessType</key>
    <string>Background</string>
    
    <key>LowPriorityIO</key>
    <true/>
    
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
